# docker/certbot/crontab
1 6,18 * * * sleep $(expr $RANDOM % 55)m; /./certbot-auto renew --quiet --no-self-upgrade